# PackageFile/debian/postinst
#/bin/sh

systemctl daemon-reload
systemctl enable AutoReport.service
systemctl restart AutoReport.service